# supercollider-arch/supervisord.ini
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:jack]
user=makepkg-user
command=/usr/bin/jackd -r -d dummy

[program:ffmpeg]
user=makepkg-user
redirect_stderr=true
command=/usr/bin/ffmpeg -f jack -i ffmpeg -c:a libmp3lame -ab 48k -ar 48000 -f hls -hls_time 1 -hls_wrap 10 -hls_allow_cache 0 -method PUT -use_localtime 1 -hls_segment_filename 'http://front-streamer/live/sc-%%Y%%m%%d-%%s.ts' http://front-streamer/live/sc.m3u8
