# supercollider-arch/supervisord.ini
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:jack]
user=makepkg-user
command=/usr/bin/jackd -r -d dummy

[program:ffmpeg]
user=makepkg-user
command=/usr/bin/ffmpeg -f jack -i ffmpeg -c:a libmp3lame -ab 48k -ar 48000 -f mpegts udp://172.17.0.1:33333
