# supercollider-arch/supervisord.ini
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:jack]
user=sc
command=/usr/bin/jackd -r -d dummy

[program:ffmpeg_egress]
user=sc
redirect_stderr=true
command=/usr/bin/ffmpeg -f jack -i ffmpeg -c:a copy -c:a libmp3lame -ab 48k -ar 48000 -f hls -hls_time 1 -hls_wrap 10 -hls_allow_cache 0 -method PUT -use_localtime 1 -hls_segment_filename 'http://front-streamer/live/sc-%%Y%%m%%d-%%s.ts' http://front-streamer/live/sc.m3u8

[program:xvfb]
user=sc
command=/usr/bin/Xvfb :0 -screen 0 640x480x8 -nolisten tcp

[program:wetty]
user=sc
redirect_stderr=true
environment=HOME="/home/sc",USER="sc"
directory=/home/sc
command=/usr/bin/node /usr/local/wetty/app.js -p 3000
